# systemd/job-agent@tailor.timer
[Unit]
Description=Run tailor hourly (offset 10m)
[Timer]
OnCalendar=*:10
Persistent=true
[Install]
WantedBy=timers.target