# systemd/job-agent@fetcher.timer
[Unit]
Description=Run fetcher hourly
[Timer]
OnCalendar=hourly
Persistent=true
[Install]
WantedBy=timers.target